# .config/profile.d/shirodeb.sh
export PATH=$PATH:"/home/shiroko/Projects/ShiroDEB/bin"
export DEB_UPLOAD_PATH=$HOME/Check
alias sd="shirodeb"
